# systemd/.config/systemd/user/alfa.service
[Unit]
Description=Radio Alfa

[Service]
ExecStart=/run/current-system/sw/bin/mplayer http://89.249.7.68/alfa

[Install]
WantedBy=multi-user.target
